# tests/mock_system/etc/systemd/system/ssh.service
[Unit]
Description=OpenBSD Secure Shell server
[Service]
Type=notify
ExecStart=/usr/sbin/sshd -D
[Install]
WantedBy=multi-user.target
